# client/crontab
* * * * * root /usr/bin/scoreboard-reporter.sh
